# modeles/postgresql-template.service
.include /lib/systemd/system/postgresql-9.6.service
#	...make your changes here...
# For more info about custom unit files, see
# http://fedoraproject.org/wiki/Systemd#How_do_I_customize_a_unit_file.2F_add_a_custom_unit_file.3F

[Unit]

[Service]

User=USER_OWNER
Group=GROUP_OWNER

# Note: avoid inserting whitespace in these Environment= lines, or you may
# break postgresql-setup.

# Location of database directory
Environment=
Environment=PGDATA=/custom/data/pgsql/9.6/NOM_INSTANCE/
Environment=PGLOG=/custom/log/pgsql/9.6/NOM_INSTANCE/
